# n8n/n8n.service
[Unit]
Description=n8n workflow automation (Docker)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=peter
WorkingDirectory=/home/peter/hierarchical-agents/n8n
EnvironmentFile=/home/peter/hierarchical-agents/n8n/.env
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
