# dashboard/atlas-watcher.service
[Unit]
Description=Event-driven atlas watcher (wakes atlas only on dashboard/delegated messages)
After=network.target agent-dashboard.service

[Service]
Type=simple
User=peter
ExecStart=/bin/bash /home/peter/hierarchical-agents/dashboard/atlas-watcher.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
